# dynv6/dynv6.timer
[Unit]
Description=dynv6 update timer

[Timer]
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
